# ai-upscale-server/plex-ai-upscale.service
[Unit]
Description=Plex Play AI Upscale Server
After=network-online.target plexmediaserver.service
Wants=network-online.target

[Service]
Type=simple
User=plex
Group=plex
EnvironmentFile=/etc/plex-ai-upscale.env
WorkingDirectory=/opt/plex-ai-upscale/realesrgan
ExecStart=/usr/bin/python3 /opt/plex-ai-upscale/server.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
